# crontab
1 * * * * /usr/local/bin/python /updater.py >> /logfile
